# server_fs/srv/client_fs/etc/local.d/050-freerdp-install.start
#!/bin/sh
#### Startup script to perform any additional tasks before system is handed over to user

# Install custom FreeRDP package
apk update --allow-untrusted
apk add /opt/freerdp-usbmod-libs-2.11.7-r1.apk --allow-untrusted --force-non-repository
apk add /opt/freerdp-usbmod-2.11.7-r1.apk --allow-untrusted --force-non-repository

# Make xfreerdp binary suid for USB access because udev rules just aren't working
chmod 4755 /usr/bin/xfreerdp

