# server_fs/srv/client_fs/etc/local.d/060-get-config.start
#!/bin/sh

# Client script for retrieving modified rdpconfig.ini if it exists on server

SERVER="$(cat /etc/next_server)"

wget_status="$(wget http://${SERVER}/alpine/client/rdpconfig.ini -O /tmp/rdpconfig.ini 2>&1 | grep 'HTTP/' | awk '{print $2}')"

if [ "${wget_status}" = "200" ]; then
        rm /etc/rdpconfig.ini
        mv /etc/rdpconfig.ini /etc
fi
