# server_fs/srv/client_fs/etc/local.d/060-get-config.start
#!/bin/sh

# Client script for retrieving modified rdpconfig.ini if it exists on server

SERVER="$(cat /etc/next_server)"

wget "http://${SERVER}/alpine/client/rdpconfig.ini" -O /etc/rdpconfig.ini

