# server_fs/srv/client_fs/etc/local.d/099-last-call.start
#!/bin/sh

# Last call script to retrieve client 'last_call.sh' from 'next_server'
# If last_call.sh doesn't exist, next step is login session

SERVER="$(cat /etc/next_server)"

wget -qO - "http://${SERVER}/alpine/client/last_call.sh"

